# systemd/user/redshift.service
[Unit]
Description=Guake Terminal

[Service]
ExecStart=/usr/bin/redshift -l '50:14'
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
